# doc/demo/docker/gentoo/etc/local.d/omnidb.start
pushd /OmniDB/OmniDB || exit 1
./omnidb-server.py --host 0.0.0.0 &
popd
